# chronoline.service
[Unit]
Description=Chronoline Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/chronoline/backend
ExecStart=/usr/bin/node simple-server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target 